# minctweet.conf
description "Minecraft Tweeter"
author  "Matt Williams"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/minecraft/minctweet
exec /home/minecraft/minctweet/minctweet.py
respawn